# system/cbdsetup-pre.target
[Unit]
Description=CloudBD Volumes (Pre)
RefuseManualStart=yes
