# system/cbdsetup.target.tmpl
[Unit]
Description=CloudBD Volumes
Documentation=man:cbdtab(5)
SourcePath=@sysconfdir@/cloudbd/cbdtab

