# system/cbdsetup@.service.tmpl
[Unit]
Description=CloudBD Disk %i
Documentation=man:cbdtab(5) man:cbdsetup-generator(8) man:cbdsetup@.service(8)
SourcePath=@sysconfdir@/cloudbd/cbdtab
DefaultDependencies=no
Wants=network-online.target remote-fs-pre.target
After=network-online.target local-fs.target cbdsetup-pre.target
Before=remote-fs-pre.target cbdsetup.target

[Service]
Type=forking
ExecStart=@sbindir@/cbddisks_start %i
ExecStop=@sbindir@/cbddisks_stop -f %i
LimitMEMLOCK=infinity
